# xremap/xremap.service
[Service]
Type=simple
ExecStart=/usr/local/bin/xremap /home/pawel/dotfiles/xremap/hjkl

[Install]
WantedBy=multi-user.target
